# files/99-atem.rules
# 99-atem.rules
# Deployed to: /etc/udev/rules.d/99-atem.rules
#
# Detects USB V4L2 capture devices and creates a stable /dev/atem_video
# symlink, then starts/stops the display service on plug/unplug.
#
# Matches any USB video capture device — primarily Blackmagic ATEM mixers
# but also webcams and other UVC capture cards for testing.
# ENV{ID_V4L_CAPABILITIES} filters to capture-capable nodes only,
# avoiding spurious matches on metadata/output nodes the same device may expose.

# On plug-in: create symlink and start display service
SUBSYSTEM=="video4linux", ACTION=="add", SUBSYSTEMS=="usb", \
    ENV{ID_V4L_CAPABILITIES}=~":capture:", \
    SYMLINK+="atem_video", \
    RUN+="/bin/systemctl start atem-display.service"

# On removal: stop display service (script loop will wait for re-plug)
SUBSYSTEM=="video4linux", ACTION=="remove", SUBSYSTEMS=="usb", \
    ENV{ID_V4L_CAPABILITIES}=~":capture:", \
    RUN+="/bin/systemctl stop atem-display.service"
