# files/ssh-keygen-firstboot.service
# ssh-keygen-firstboot.service
# Deployed to: /etc/systemd/system/ssh-keygen-firstboot.service
#
# Generates SSH host keys on first boot if they are absent.
# Pi OS Trixie removes host keys in the base image; this ensures they exist
# before sshd starts, without relying on distro-specific regeneration scripts.

[Unit]
Description=Generate SSH host keys on first boot
Before=ssh.service
ConditionPathExists=!/etc/ssh/ssh_host_rsa_key

[Service]
Type=oneshot
ExecStart=/usr/bin/ssh-keygen -A
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
